# lb2002_dtb_research/armbian_files_cm311_hmi/etc/systemd/system/azs-hmi-cm311.service.d/20-hdmi-bpc8.conf
[Service]
ExecStartPre=/usr/local/sbin/cm311-set-hdmi-bpc8 35
